# templates/default/upstart-content-server.conf.erb
# description "Start the fxa-content-server"

start on runlevel [2345]
stop on runlevel [^2345]

console log
chdir <%= node['mozilla-firefox-accounts']['content-server']['path'] %>
setuid fxcloud
setgid fxcloud

env NODE_ENV=prod 
respawn
respawn limit 20 5

exec npm start

